# playbooks/files/update-motd.d/90-updates-available
#!/bin/bash
UPDATE_COUNT=$(apt list --upgradable|grep -v Listing...|wc -l)
echo "${UPDATE_COUNT} updates can be applied immediately."
